# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/busybox-syslog.list
/.
/etc
/etc/init.d
/etc/init.d/syslog
/etc/syslog-startup.conf
/etc/syslog.conf
