# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/dropbear.list
/.
/etc
/etc/default
/etc/default/dropbear
/etc/dropbear
/etc/init.d
/etc/init.d/dropbear
/usr
/usr/bin
/usr/sbin
/usr/sbin/dropbearmulti
/var
/usr/bin/dbclient
/usr/sbin/dropbear
/usr/sbin/dropbearconvert
/usr/sbin/dropbearkey
