# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/initscripts.list
/.
/etc
/etc/default
/etc/default/devpts
/etc/default/volatiles
/etc/default/volatiles/00_core
/etc/init.d
/etc/init.d/banner.sh
/etc/init.d/bootmisc.sh
/etc/init.d/checkroot.sh
/etc/init.d/devpts.sh
/etc/init.d/dmesg.sh
/etc/init.d/halt
/etc/init.d/hostname.sh
/etc/init.d/mountall.sh
/etc/init.d/mountnfs.sh
/etc/init.d/populate-volatile.sh
/etc/init.d/read-only-rootfs-hook.sh
/etc/init.d/reboot
/etc/init.d/rmnologin.sh
/etc/init.d/save-rtc.sh
/etc/init.d/sendsigs
/etc/init.d/single
/etc/init.d/sysfs.sh
/etc/init.d/umountfs
/etc/init.d/umountnfs.sh
/etc/init.d/urandom
/etc/logrotate-dmesg.conf
/etc/rc0.d
/etc/rc1.d
/etc/rc2.d
/etc/rc3.d
/etc/rc4.d
/etc/rc5.d
/etc/rc6.d
/etc/rcS.d
/var
/var/lib
/var/lib/urandom
/etc/rc0.d/S20sendsigs
/etc/rc0.d/S25save-rtc.sh
/etc/rc0.d/S31umountnfs.sh
/etc/rc0.d/S38urandom
/etc/rc0.d/S40umountfs
/etc/rc0.d/S90halt
/etc/rc1.d/S31umountnfs.sh
/etc/rc2.d/S15mountnfs.sh
/etc/rc2.d/S99rmnologin.sh
/etc/rc3.d/S15mountnfs.sh
/etc/rc3.d/S99rmnologin.sh
/etc/rc4.d/S15mountnfs.sh
/etc/rc4.d/S99rmnologin.sh
/etc/rc5.d/S15mountnfs.sh
/etc/rc5.d/S99rmnologin.sh
/etc/rc6.d/S20sendsigs
/etc/rc6.d/S25save-rtc.sh
/etc/rc6.d/S31umountnfs.sh
/etc/rc6.d/S38urandom
/etc/rc6.d/S40umountfs
/etc/rc6.d/S90reboot
/etc/rcS.d/S02banner.sh
/etc/rcS.d/S02sysfs.sh
/etc/rcS.d/S03mountall.sh
/etc/rcS.d/S06checkroot.sh
/etc/rcS.d/S06devpts.sh
/etc/rcS.d/S29read-only-rootfs-hook.sh
/etc/rcS.d/S36bootmisc.sh
/etc/rcS.d/S37populate-volatile.sh
/etc/rcS.d/S38dmesg.sh
/etc/rcS.d/S38urandom
/etc/rcS.d/S39hostname.sh
