# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/initscripts-functions.list
/.
/etc
/etc/init.d
/etc/init.d/functions
