# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/init-ifupdown.list
/.
/etc
/etc/init.d
/etc/init.d/networking
/etc/network
/etc/network/if-down.d
/etc/network/if-post-down.d
/etc/network/if-pre-up.d
/etc/network/if-pre-up.d/nfsroot
/etc/network/if-up.d
/etc/network/interfaces
/etc/network/nm-disabled-eth0
