# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/modutils-initscripts.list
/.
/etc
/etc/init.d
/etc/init.d/modutils.sh
