# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/modutils-initscripts.postinst
#!/bin/sh
set -e
	if false; then
		if [ -n "$D" ]; then
			OPTS="--root=$D"
		fi
		systemctl $OPTS mask modutils.service
	fi
if true && type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-r $D"
	else
		OPT="-s"
	fi
	update-rc.d $OPT modutils.sh start 05 S .
fi
