# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/run-postinsts.list
/.
/etc
/etc/init.d
/etc/init.d/run-postinsts
/usr
/usr/sbin
/usr/sbin/run-postinsts
