# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/sysvinit.list
/.
/bin
/bin/mountpoint.sysvinit
/etc
/etc/default
/etc/default/rcS
/etc/default/volatiles
/etc/default/volatiles/01_bootlogd
/etc/init.d
/etc/init.d/bootlogd
/etc/init.d/rc
/etc/init.d/rcS
/etc/rc0.d
/etc/rc1.d
/etc/rc2.d
/etc/rc3.d
/etc/rc4.d
/etc/rc5.d
/etc/rc6.d
/etc/rcS.d
/sbin
/sbin/bootlogd
/sbin/fstab-decode
/sbin/halt.sysvinit
/sbin/init.sysvinit
/sbin/runlevel.sysvinit
/sbin/shutdown.sysvinit
/usr
/usr/bin
/usr/bin/last.sysvinit
/usr/bin/mesg.sysvinit
/usr/bin/utmpdump.sysvinit
/usr/bin/wall.sysvinit
/etc/init.d/stop-bootlogd
/etc/rc2.d/S99stop-bootlogd
/etc/rc3.d/S99stop-bootlogd
/etc/rc4.d/S99stop-bootlogd
/etc/rc5.d/S99stop-bootlogd
/etc/rcS.d/S07bootlogd
/sbin/poweroff.sysvinit
/sbin/reboot.sysvinit
/sbin/telinit
/usr/bin/lastb.sysvinit
