# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xendomains.list
/.
/etc
/etc/default
/etc/default/xendomains
/etc/init.d
/etc/init.d/xendomains
/etc/sysconfig
/usr
/usr/lib
/usr/lib/xen
/usr/lib/xen/bin
/usr/lib/xen/bin/xendomains
/etc/sysconfig/xendomains
