# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xen-watchdog.postrm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
if true && type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-f -r $D"
	else
		OPT="-f"
	fi
	update-rc.d $OPT xen-watchdog remove
fi
