# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xencommons.preinst
#!/bin/sh
if true && [ -z "$D" -a -f "/etc/init.d/xencommons" ]; then
	/etc/init.d/xencommons stop || :
fi
if true && type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-f -r $D"
	else
		OPT="-f"
	fi
	update-rc.d $OPT xencommons remove
fi
