# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xencommons.list
/.
/etc
/etc/default
/etc/default/xencommons
/etc/init.d
/etc/init.d/xencommons
/etc/xen
/etc/xen/scripts
/etc/xen/scripts/launch-xenstore
