# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xen-watchdog.list
/.
/etc
/etc/init.d
/etc/init.d/xen-watchdog
/usr
/usr/sbin
/usr/sbin/xenwatchdogd
