# Software_Driver/Co-Sim/linux/linux/targetfs.arm64/var/lib/dpkg/info/xen-xen-watchdog.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
# Begin section update-rc.d
if true && [ -z "$D" -a -x "/etc/init.d/xen-watchdog" ]; then
	/etc/init.d/xen-watchdog stop || :
fi
# End section update-rc.d
